# scripts/admin/backup/btrbk/templates/btrbk-template.timer
[Unit]
Description={{DESCRIPTION}} - Timer Unit

[Timer]
OnCalendar={{FREQUENCY}}
Persistent=true

[Install]
WantedBy=timers.target
